# cron/crontab
30 *    * * *   root    /opt/rddmarc/rddmarc /home/<===dmarc-mail-recipient===>/Maildir/dmarc/new/* &> /dev/null#

